# crash_reporter/init/warn-collector.conf
# Copyright (c) 2013 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Runs a daemon which collects and reports kernel warnings"
author      "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services
respawn

exec warn_collector
